# backups/2026-03-23/system/openclaw.service
# /etc/systemd/system/openclaw.service
[Unit]
Description=OpenClaw Gateway
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw
Environment="HOME=/root"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="OPENCLAW_CONFIG=/root/.openclaw"
Environment="DBUS_SESSION_BUS_ADDRESS=/dev/null"
ExecStart=/bin/bash -c 'cd /root/.openclaw && exec openclaw gateway run'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
